# templates/automatic-upgrades.timer.erb
[Unit]
Description=Schedule pacman upgrades

[Timer]
OnBootSec=<%= @bootdelay %>
OnUnitActiveSec=<%= @frequency %>

[Install]
WantedBy=multi-user.target

