# templates/automatic-upgrades.service.erb
[Unit]
Description=Run pacman upgrades

[Service]
Type=simple
ExecStart=/usr/bin/pacman -Syyu --noconfirm

