# files/pacman-cache-clean.service
[Unit]
Description=Run pacman cache clean

[Service]
Type=simple
ExecStart=/usr/bin/pacman -Sc --noconfirm

