# files/pacman-cache-clean.timer
[Unit]
Description=Schedule pacman cache clean

[Timer]
OnCalendar=weekly

[Install]
WantedBy=multi-user.target

